# fuel-plugin/debian/qemu-system-common.qemu-kvm.upstart
# qemu-kvm

description "KVM"
author "Dustin Kirkland <kirkland@canonical.com>"

start on runlevel [2345]

pre-start exec /usr/share/qemu/init/qemu-kvm-init start
